# systemd/user/process_pillz-restarter.path
[Unit]
Description=Watch changes to process_pillz's config

[Path]
PathChanged=%h/.config/process_pillz.yaml

[Install]
WantedBy=default.target

