# systemd/user/process_pillz.service
[Unit]
Description=Process Pillz
After=network.target

[Service]
ExecStart=process_pillz
Type=simple
Restart=always
RestartSec=1
SuccessExitStatus=42
#StandardOutput=journal
#StandardError=journal

# Security options
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

# Only allow necessary capabilities
CapabilityBoundingSet=CAP_SYS_NICE CAP_SYSLOG
AmbientCapabilities=CAP_SYS_NICE

# Restrict filesystem access
PrivateTmp=yes
ReadOnlyPaths=/
ReadWritePaths=/proc

[Install]
WantedBy=default.target
